# cfy_manager/components/mgmtworker/config/supervisord.conf
[program:cloudify-mgmtworker]
command=/opt/mgmtworker/env/bin/python -m mgmtworker.worker --queue "cloudify.management" --max-workers {{ mgmtworker.max_workers }} --hooks-queue "cloudify-hooks" --cluster-service-queue {{ mgmtworker.cluster_service_queue | default('cluster-queue') }}
environment=MGMTWORKER_HOME="{{ mgmtworker.home_dir }}",VIRTUALENV_DIR="{{ mgmtworker.home_dir }}/env",AGENT_WORK_DIR="{{ mgmtworker.home_dir }}/work",AGENT_LOG_DIR="{{ mgmtworker.log_dir }}",AGENT_LOG_LEVEL="{{ mgmtworker.log_level }}",MANAGER_NAME="{{ manager.hostname }}",REST_HOST="{{ networks.default }}",REST_PORT="{{ constants.internal_rest_port }}",LOCAL_REST_CERT_FILE="{{ constants.ca_cert_path }}",BROKER_SSL_CERT_PATH="{{ rabbitmq.ca_path }}"
autostart=false
