# cfy_manager/components/nginx/config/supervisord.conf
[program:nginx]
autorestart=true
minfds={{ nginx.max_open_fds }}
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf -g 'daemon off;'
