# cfy_manager/components/manager_ip_setter/config/supervisord/cloudify-manager-ip-starter.conf
[program:cloudify-manager-ip-starter]
autostart=false
command=/opt/cloudify/manager-ip-setter/manager_ip_starter.sh
user=root
group=root
