# cfy_manager/components/manager_ip_setter/config/supervisord/cloudify-manager-ip-setter.conf
[program:cloudify-manager-ip-setter]
autostart=false
command=/opt/cloudify/manager-ip-setter/manager-ip-setter.sh
user=root
group=root
