# cfy_manager/components/postgresql_server/supervisord.conf
[program:cloudify-postgresql-server]
command=/usr/pgsql-9.5/bin/postgres -D /var/lib/pgsql/9.5/data/
autostart=false
autorestart=false
startsecs=0
environment=PGDATA="/var/lib/pgsql/9.5/data/"