# cfy_manager/components/status_reporter/supervisord.conf
[program:cloudify-status-reporter]
autostart=false
command=/opt/status-reporter/env/bin/cloudify_{{ reporter_type }} {{ extra_config_flags }}
