# cfy_manager/components/syncthing/supervisord.conf
[program:cloudify-syncthing]
command=/opt/syncthing/syncthing -no-restart
environment=STNORESTART="1",STNOUPGRADE="1",HOME="/opt/syncthing"
