# cfy_manager/components/nginx/config/supervisord/wait_on_restart.conf
[program:wait_on_restart]
autostart=false
startsecs=0
command=/etc/cloudify/wait_on_restart.sh