# cfy_manager/components/nginx/config/supervisord/wait_on_restart.conf
[program:wait_on_restart]
autostart=false
command=/opt/cloudify/wait_on_restart.sh
