# cfy_manager/components/restservice/config/supervisord.conf
[program:cloudify-restservice]
user=cfyuser
group=cfyuser
autostart=false
command=/bin/sh -c '/opt/manager/env/bin/gunicorn \
    --pid /run/cloudify-restservice/pid \
    --chdir / \
    --workers 8 \
    --max-requests 1000 \
    --bind 127.0.0.1:8100 \
    --log-level debug \
    --timeout 300 manager_rest.wsgi:app \
    --error-logfile /var/log/cloudify/rest/error.log \
    --access-logfile /var/log/cloudify/rest/gunicorn-access.log'
environment=HOME="/etc/cloudify",USER="cfyuser",MANAGER_REST_CONFIG_PATH="{{ restservice.home_dir }}/cloudify-rest.conf",MANAGER_REST_SECURITY_CONFIG_PATH="{{ restservice.home_dir }}/rest-security.conf",MANAGER_REST_AUTHORIZATION_CONFIG_PATH="{{ restservice.home_dir }}/authorization.conf",REST_PORT="{{ restservice.port }}"
