# cfy_manager/components/rsyslog/config/supervisord/rsyslog.conf
[program:rsyslog]
startsecs=0
umask=0066
command=/opt/cloudify/syslog_wrapper_script.sh