# cfy_manager/components/rsyslog/config/supervisord/rsyslog.conf
[program:rsyslog]
stdout_logfile = /var/log/cloudify/rsyslog.log
redirect_stderr = true
startsecs=0
umask=0066
command=/opt/cloudify/syslog_wrapper_script.sh
