# cfy_manager/components/stage/config/cloudify-stage.service
[Unit]
Description=Cloudify Console Service
Wants=network-online.target
After=network-online.target

[Service]
Type=forking
LimitNOFILE=infinity
LimitNPROC=infinity
LimitCORE=infinity
Restart=on-failure
RestartSec=10
PIDFile=/opt/cloudify-stage/.pm2/pm2.pid
EnvironmentFile=/etc/sysconfig/cloudify-stage
User=cfyuser
Group=cfyuser
WorkingDirectory=/opt/cloudify-stage/backend
ExecStart=/opt/nodejs/bin/npm run start
ExecStop=/opt/nodejs/bin/npm run stop
ExecReload=/opt/nodejs/bin/npm run restart
UMask=0002

[Install]
WantedBy=multi-user.target
