# cfy_manager/components/restservice/config/supervisord/configure-rest-runtime-dir.conf
[program:configure-rest-runtime-dir]
command=/opt/cloudify/configure-rest-runtime-dir.sh
startsecs=0
