# cfy_manager/components/postgresql_server/config/supervisord/patroni_startup_check.conf
[program:patroni_startup_check]
autostart=false
command=/opt/patroni/bin/patroni_startup_check