# packaging/supervisord/rsyslog.conf
[program:rsyslog]
autostart=true
startsecs=0
umask=0066
environment=SYSLOGD_OPTIONS=""
command=/usr/sbin/rsyslogd -n $SYSLOGD_OPTIONS