# packaging/supervisord/rsyslog.conf
[program:rsyslog]
autostart=true
startsecs=0
umask=0066
command=/usr/sbin/rsyslogd -n