# packaging/supervisord/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
environment=SSH_USE_STRONG_RNG=0
