# packaging/supervisord/sshd.conf
[program:sshd]
autostart=true
command=/usr/sbin/sshd
